# systemd/nvidiasmi_exporter.service
[Unit]
Description=NvidiaSmi_Exporter

[Service]
ExecStart=/opt/prometheus/nvidiasmi_exporter
Restart=always

[Install]
WantedBy=multi-user.target
